# deploy/roles/infra_server/templates/chem-infra-api.service
[Unit]
Description=Chemhacktica Infra API Server
After=network.target

[Service]
ExecStart=/home/askcos/infra/start-api.sh
Type=exec
Restart=always
User=askcos
WorkingDirectory=/home/askcos/infra

[Install]
WantedBy=default.target
RequiredBy=network.target
